# intern/boot_time.service
[Unit]
Description=Boot Time Script

[Service]
ExecStart=/bin/bash ~/intern_DGIST/boot_time.sh

[Install]
WantedBy=default.target
